# deploy/com.jarvis.credential-manager.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.jarvis.credential-manager</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/jarvis/jarvis-credential-manager/deploy/launch.sh</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>WorkingDirectory</key>
    <string>/Users/jarvis/jarvis-credential-manager</string>

    <key>StandardOutPath</key>
    <string>/Users/jarvis/Library/Logs/jarvis-credential-manager.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/jarvis/Library/Logs/jarvis-credential-manager.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PORT</key>
        <string>3847</string>
        <key>REQUEST_TIMEOUT_MS</key>
        <string>600000</string>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
